# debian/cryptpilot-crypt.prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop cryptpilot.service || true
        systemctl disable cryptpilot.service || true
    fi
fi
#DEBHELPER#
exit 0
