# dist/cryptpilot.service
[Unit]
Description=Auto open cryptpilot volumes
DefaultDependencies=no
After=network-online.target
Before=remote-fs-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/cryptpilot systemd-service

[Install]
WantedBy=remote-fs-pre.target
