# dist/systemd/cryptpilot.service
[Unit]
Description=Cryptpilot Volumes Auto Open Service
Requires=network-online.target
After=network-online.target
After=cryptpilot-fde-before-sysroot.service
Wants=attestation-agent.service
After=attestation-agent.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/cryptpilot-crypt boot-service --stage system-volumes-auto-open
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
